# ImageServer.service
[Unit]
Description=ImageServer
After=network.target

[Service]
Type=simple
RemainAfterExit=true
WorkingDirectory=/etc/server
ExecStart=/etc/server/ImageServer
PrivateTmp=false

[Install]
WantedBy=multi-user.target
